# systemd/webcam-shoot.service
[Unit]
Description=Webcam shoot and store

[Service]
ExecStart=/home/knockside/webcam/shoot

[Install]
WantedBy=default.target
