# systemd/webcam-shoot.timer
[Unit]
Description=Webcam shoot and store timer

[Timer]
OnCalendar=*-*-* *:0/5:00

[Install]
WantedBy=default.target
